# Roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=node_exporter service
After=network.target

[Service]
ExecStart=/usr/local/bin/node_exporter
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target